# systemd/user/check-battery-user.timer
[Unit]
Description=Timer to check battery status

[Timer]
OnActiveSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
