# systemd/user/check-battery-user.service
[Unit]
Description=check-battery

[Service]
ExecStart=/home/mike/.config/systemd/user/check-battery.sh

[Install]
WantedBy=default.target
