# systemd/user/doc_sync.timer
[Unit]
Description=Timer to sync docs

[Timer]
OnActiveSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
