# systemd/user/nfs_share.service
[Unit]
Description=mount_nfs
After=network-online.target
Wants=network-online.target


[Service]
Type=oneshot
ExecStartPre=/bin/sh -c 'until ping -c1 1.1.1.1 2>/dev/null; do sleep 1; done;'
ExecStart=/home/mike/.config/systemd/user/nfs_share.sh

[Install]
WantedBy=default.target

